# deploy/systemd/movie_recommender-scraper.service
[Unit]
Description=Movie Recommender Scraper
After=network.target

[Service]
Type=simple
# Optionally set User to a system user that has permissions to the project directory.
# If left unset the service will run as root. Adjust before enabling on your host.
WorkingDirectory=/home/abhishek/movie_recommender
Environment=PATH=/home/abhishek/movie_recommender/.venv/bin:/usr/bin:/bin
Environment=TMDB_API_KEY=
ExecStart=/home/abhishek/movie_recommender/.venv/bin/python -m src.scraper --n 200 --append --interval 3600
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
